# systemd/user/mpd.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/mpd --systemd /home/_/.mpd/mpd.conf
